# lib/udev/rules.d/85-tdp.rules
ACTION=="change", SUBSYSTEM=="power_supply", RUN+="/usr/local/bin/tdp auto"
